# supervisord.conf
[program:abbott]
command = PLATFORM=dotcloud npm start
directory = /home/dotcloud/current