# supervisord.conf
[program:abbott]
command = coffee bin/cli.coffee --platform dotcloud
directory = /home/dotcloud/current