# supervisord.conf
[program:abbott]
command = npm start
environment = PLATFORM=dotcloud
directory = /home/dotcloud/current